# images/ubuntu16.04/supervisord.conf
[supervisord]
nodaemon=true

[program:creator]
command=mkdir -p /var/run/sshd

[program:sshd]
command=/usr/sbin/sshd -D
